# apolo.service
[Unit]
Description=Redis container
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker start -a apollocontainer
ExecStop=/usr/bin/docker stop -t 2 apollocontainer

[Install]
WantedBy=default.target