# test.service
[Unit]
Description=Redis container
After=multi-user.target

[Service]
ExecStart=/usr/bin/python3 /home/Argos/apollo/lightsTest.py

[Install]
WantedBy=default.target